# workdir/debian/ethercalc.cron
# Backup of Ethercalc.
35 3 * * * root [ -x /usr/bin/ethercalc-backup ] && /usr/bin/ethercalc-backup


